# setup/znc.conf
[program:znc]
user = nobody
command = /usr/bin/znc -d /config -f
umask = 002
